# home/.config/profile.d/00-composer.sh
#!/usr/bin/env sh

export COMPOSER_IGNORE_PLATFORM_REQS=1
