# home/.config/profile.d/00-docker.sh
#      _            _
#   __| | ___   ___| | _____ _ __
#  / _` |/ _ \ / __| |/ / _ \ '__|
# | (_| | (_) | (__|   <  __/ |
#  \__,_|\___/ \___|_|\_\___|_|

export DOCKER_CONFIG="$XDG_CONFIG_HOME"/docker
