# home/.config/profile.d/00-go.sh
#!/usr/bin/env sh
#   __ _  ___
#  / _` |/ _ \
# | (_| | (_) |
#  \__, |\___/
#  |___/

export GOPATH="$XDG_DATA_HOME"/go

path_prepend "$GOPATH"/bin
