# home/.config/profile.d/00-wine.sh
#           _
# __      _(_)_ __   ___
# \ \ /\ / / | '_ \ / _ \
#  \ V  V /| | | | |  __/
#   \_/\_/ |_|_| |_|\___|

export WINEPREFIX="$XDG_DATA_HOME"/wineprefix
