# home/.config/profile.d/99-editor.sh
#!/usr/bin/env sh

export VISUAL=nano
export EDITOR="$VISUAL"
