# home/.config/profile.d/99-libgl.sh
#!/usr/bin/env sh

export LIBGL_DRI3_DISABLE=1
